# etc/cron.daily/time_sync.sh
#!/bin/bash
rdate -s time.bora.net && date && clock -r && clock -w > /dev/null 2>&1
